# reverb-supervisor.conf
# --------------------------------------------------------------------------
# KasirKite Supervisor WebSocket Configuration (reverb-supervisor.conf)
# --------------------------------------------------------------------------
# Salin berkas ini ke server GCP Anda di path:
# /etc/supervisor/conf.d/reverb-websocket.conf
# 
# Setelah disalin, jalankan perintah berikut untuk mengaktifkan daemon:
# sudo supervisorctl reread
# sudo supervisorctl update
# sudo supervisorctl start reverb-websocket
# --------------------------------------------------------------------------

[program:reverb-websocket]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/kasir-app/artisan reverb:start --host=127.0.0.1 --port=8080
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/kasir-app/storage/logs/reverb.log
stopwaitsecs=3600
